# services/captive-portal-proxy.service
[Unit]
Description=Captive Portal Detection Proxy
Documentation=https://github.com/ossuary/ossuary-pi
After=NetworkManager.service wifi-connect.service
Wants=wifi-connect.service
# Give wifi-connect time to start on port 8080
After=wifi-connect.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/ossuary/scripts/captive-portal-proxy.py
Restart=always
RestartSec=5
User=root

# Allow binding to privileged port 80
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=captive-portal-proxy

# Resource limits
LimitNOFILE=65535
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
